# supervisor/nginx.ini
[program:nginx]
command=/usr/sbin/nginx
user=root
redirect_stderr=true
stdout_logfile=/var/log/shared/supervisor.nginx.log
