# supervisor/phpfpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm -F
user=root
redirect_stderr=true
stdout_logfile=/var/log/shared/supervisor.phpfpm.log
